# monitor/lib/udev/rules.d/99-custom-monitor.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/krajeswaran/.Xauthority", RUN+="/sbin/hotplug.sh"
