# autostart/.config/autostart/albert.desktop
[Desktop Entry]
Comment=
Encoding=UTF-8
Exec=albert
Hidden=false
Name=albert
StartupNotify=false
Terminal=false
Type=Application
Version=0.9.4

